# etc/systemd/system/snap-snapd-25935.mount
[Unit]
Description=Mount unit for snapd, revision 25935
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/snapd_25935.snap
Where=/snap/snapd/25935
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
